# deploy/systemd/openclaw-zabbix-mcp.service
[Unit]
Description=Zabbix MCP API (FastAPI)
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
WorkingDirectory=/opt/404_donkey_not_found/apps/openclaw_zabbix_mcp
ExecStart=/opt/404_donkey_not_found/apps/openclaw_zabbix_mcp/venv/bin/python3 -m uvicorn app:app --host 127.0.0.1 --port 8888
Restart=always
RestartSec=3
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
